# bot.service
[Unit]
Description=Personal Assistant Telegram Bot
After=network.target

[Service]
User=root
WorkingDirectory=/opt/personal-assistant-bot
ExecStart=/opt/personal-assistant-bot/venv/bin/python main.py
Restart=always
RestartSec=5
# Adjust these paths depending on where you clone the repo on your Debian server

[Install]
WantedBy=multi-user.target
